# etc/systemd/system/supermicro-fan-control.service
[Unit]
Description=Supermicro Fan Controller Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash -c "cd /opt/supermicro-fan-control && source venv/bin/activate && python3 app/supermicro-fan-control.py"
Restart=on-failure

[Install]
WantedBy=multi-user.target
